# scripts/ops/db-load-observability.service
[Unit]
Description=Longhouse SQLite load observability sampler
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
EnvironmentFile=/etc/longhouse-db-observability.env
ExecStart=/usr/bin/python3 /usr/local/lib/longhouse/db-load-observability.py sample --runtime-container ${OBSERVABILITY_RUNTIME_CONTAINER} --containers ${OBSERVABILITY_RESOURCE_CONTAINERS}
